# files-giot/etc/uci-defaults/99_wireless.erb
uci -m import wireless < /var/tmp/wireless.uci
rm /var/tmp/wireless.uci
uci set wireless.radio1.disabled='0'
uci commit wireless
wifi
